# build/supervisor/conf.d/cron.conf
[program:crond]
nodaemon = true
process_name=%(program_name)s_%(process_num)02d
command=busybox crond -l 2 -L /dev/stderr -f
numprocs=1
autostart=true
autorestart=false
redirect_stderr=true
stopwaitsecs=120
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
